# contrib/feishu-bridge@.service
[Unit]
Description=Feishu Bridge (%i)
After=network-online.target

[Service]
Type=simple
ExecStart=%h/.local/bin/feishu-bridge --bot %i
Restart=always
RestartSec=2

# Environment (adjust path as needed)
EnvironmentFile=-%h/.config/feishu-bridge/.env

# Security hardening
MemoryMax=2G
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=%h/.local %h/.config %h/.openclaw %h/.claude /tmp
StandardOutput=journal
StandardError=journal

# Prevent infinite restart loops
StartLimitIntervalSec=120
StartLimitBurst=8

[Install]
WantedBy=default.target
